# .local/runit/sv/warp/run
#!/bin/sh

exec 2>&1

exec /opt/warp/bin/warp-svc
