# .local/runit/sv/syncthing/run
#!/bin/sh

export USERNAME=mustaqim
export HOME="/home/$USERNAME"
export SYNCTHING="/usr/bin/syncthing"

exec 2>&1
exec chpst -u "$USERNAME" "$SYNCTHING" -logflags 0
